# .linux/systemctl/homepage.service
[Unit]
Description=Homepage Service
After=network.target

[Service]
Environment=ENV=prod
ExecStart=/home/ubuntu/app/homepage
WorkingDirectory=/home/ubuntu/app
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target

